# etc/systemd/system/prune-sent-mail.service.example
[Unit]
Description=Prune Sent Mailbox

[Service]
ExecStart=/usr/bin/python3 /path/to/prune.py '[Gmail]/All Mail' 60
